# martyrs-archive.service
[Unit]
Description=Martyrs Archive Node.js Application
Documentation=https://github.com/martyrs-archive
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/apps/Martyrs-Docments
Environment=NODE_ENV=production
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PORT=5000
Environment=CLIENT_URL=http://localhost:3000
Environment=DB_HOST=localhost
Environment=DB_USER=api
Environment=DB_PASSWORD=7e9af97ec189f5448a809ca7b3f7e7c0
Environment=DB_NAME=martyrs_archive
Environment=JWT_SECRET=your-super-secret-jwt-key-change-this-in-production
Environment=MAX_FILE_SIZE=5242880
Environment=UPLOAD_PATH=./uploads
Environment=CORS_ORIGIN=http://localhost:3000
Environment=RATE_LIMIT_WINDOW=900000
Environment=RATE_LIMIT_MAX=100
Environment=SMTP_HOST=smtp.gmail.com
Environment=SMTP_PORT=587
Environment=SMTP_USER=your-email@gmail.com
Environment=SMTP_PASS=your-app-password
Environment=GOOGLE_MAPS_API_KEY=your-google-maps-api-key
Environment=BACKUP_PATH=./backups
Environment=BACKUP_RETENTION_DAYS=30
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=10
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
SyslogIdentifier=martyrs-archive

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/apps/Martyrs-Docments
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
